# debian/postinst
#!/bin/bash -e

echo "Enabling ZIPGATEWAY service"
systemctl enable zipgateway

